# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:server]
command=node server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:log_consumer]
; Flushes Redis Stream events -> /app/logs/cache_events.jsonl every 5 min.
; This file is on the shared logs_data volume, so the ml-service can read it
; to run benchmark.py and retrain the model — all automatically.
command=node scripts/logConsumer.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
; Wait for server.js to connect to Redis before starting the consumer
startsecs=5
